# sample/atticmatic.cron
# You can drop this file into /etc/cron.d/ to run atticmatic nightly.

0 3 * * * root /usr/local/bin/atticmatic
