# sample/systemd/borgmatic.service
[Unit]
Description=borgmatic backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/borgmatic
