# sample/systemd/borgmatic.service
[Unit]
Description=borgmatic backup

[Service]
Type=oneshot
ExecStart=/root/.local/bin/borgmatic
LogRateLimitIntervalSec=0
